# profile.d/path.sh
#! /bin/bash

PATH="${PATH:+$PATH:}/devops-cli-azenvs/bin:/root/.local/bin"
export PATH
