# profile.d/path.sh
#! /bin/bash

PATH="${PATH:+$PATH:}/devops-cli-azenvs/bin"
export PATH
